# deploy/systemd.service
[Unit]
Description=Guardian detector service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/guardian
Environment=NODE_ENV=production
EnvironmentFile=-/etc/default/guardian
ExecStartPre=/usr/bin/env pnpm exec tsx src/cli.ts --health
ExecStart=/usr/bin/env pnpm exec tsx src/cli.ts daemon start
ExecStop=/usr/bin/env pnpm exec tsx src/cli.ts daemon stop
ExecStopPost=/usr/bin/env pnpm exec tsx src/cli.ts daemon hooks --reason systemd-stop --signal SIGTERM
ExecReload=/usr/bin/env pnpm exec tsx src/cli.ts daemon status --json
ExecStopPost=/usr/bin/env pnpm exec tsx scripts/db-maintenance.ts
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
